# supervisord.conf
[supervisord]
nodaemon=true

[program:go-gin-example]
command=/app/gin-blog-api
