# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:gin-blog-api]
command=/working/gin-blog-api
